# flutter_app/chatbot_server/ops/systemd/django-chatbot-cache-clear.timer
[Unit]
Description=Run Django Chatbot cache clear daily

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
